# packages/dci-umb/dci-umb-0.3.0.tar.gz/dci-umb-0.3.0/systemd/dci-umb.service
[Service]
Type=simple
EnvironmentFile=-/etc/dci-umb/config
ExecStart=/usr/bin/dci-umb
Restart=always
RestartSec=5s