# packages/pcrunner/pcrunner-0.4.8.tar.gz/pcrunner-0.4.8/systemd/pcrunner.service
[Unit]
Description=Pcrunner (Passive Check Runner)
Requires=network.target
After=network.target

[Service]
ExecStart=/usr/bin/pcrunner -a start
PIDFile=/var/run/pcrunner.pid
Restart=on-failure

[Install]
WantedBy=multi-user.target
